# nodes/brain_18/systemd/signal_fusion_notify.service
[Unit]
Description=QuantForce Apex — Signal Fusion (Notify Mode)
After=network.target postgresql.service

[Service]
Type=simple
User=heng
WorkingDirectory=/home/heng/quantforce-apex
Environment=PYTHONPATH=/home/heng/quantforce-apex
Environment=QF_PG_HOST=192.168.0.18
Environment=QF_PG_PASS=newpassword123
ExecStart=/usr/bin/python3 /home/heng/quantforce-apex/nodes/brain_18/signal_fusion_notify.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
